# orchestration/scheduler/crontab
# Supercronic schedule for the app container
* * * * * /usr/local/bin/python /app/scripts/run.py
