# orchestration/scheduler/crontab
# Supercronic schedule for the app container
* * * * * cd /app && /usr/local/bin/python -m scripts.run
